# oslokommune/jobs/crontabs
10 6,11,13,17 * * * /var/kartoteket/source/systemoversikt/oslokommune/jobs/backup_script_db >> /var/kartoteket/jobs/backup_log.txt
0 2 * * * /var/kartoteket/source/systemoversikt/oslokommune/jobs/daily_jobs >> /var/kartoteket/jobs/job_log.txt

