# oslokommune/jobs/crontabs
10 6,11,13,17 * * * /var/kartoteket/source/backup/backup_script_db >> /var/kartoteket/source/jobs/backup_log.txt
0 2 * * * /var/kartoteket/source/jobs/daily_jobs >> /var/kartoteket/source/jobs/job_log.txt

